# deployments/ansible/roles/common/files/jakeblog.service
[Unit]
Description=JakeBlog Application Service
After=docker.service
Requires=docker.service

[Service]
Type=simple
User=jakeblog
Group=jakeblog
WorkingDirectory=/opt/jakeblog
ExecStart=/usr/bin/docker compose up
ExecStop=/usr/bin/docker compose down
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
